# deploy/etc/s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv execlineb

cd /app
exec python onstock/__main__.py qcluster
