# deploy/etc/s6-overlay/s6-rc.d/server/run
#!/command/with-contenv execlineb

cd /app
exec python onstock/__main__.py litestream replicate -exec "python onstock/__main__.py"
